# danilevskii-1/etc/crontabs/root
59 4 * * * sleep 61 && touch /etc/banner && jffs2reset -y && reboot
0 * * * * wget http://ipv4.tunnelbroker.net/nic/update?username=Slenderchat&password=iG-YT3PEEb9hMvM-&hostname=546539&myip=$(ifconfig eth0.2 | sed -nE 's/.*inet addr:(.*) Bcast:.*/\1/gp')
