# crontab
59 4 * * * sleep 61 && touch /etc/banner && reboot
