# crontab
0 2 * * * sleep 70 && touch /etc/banner && jffs2reset && reboot
